# files/concrete5.cron
@reboot root /usr/local/bin/chef-solo -c /var/chef/solo.rb -j /var/chef/nodes/ec2.json & rm -f /etc/cron.d/concrete5
